# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# PostgreSQL Database - Priority 10 (starts first)
[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/14/main -c config_file=/etc/postgresql/14/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_error.log
environment=PGDATA="/var/lib/postgresql/14/main"

# Redis Server - Priority 20 (starts after PostgreSQL)
[program:redis]
command=redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_error.log

# Database Migration - Priority 30 (runs once after DB is ready)
[program:migrate]
command=python /app/core-svc/migrations/run_migrations.py
directory=/app
user=root
autostart=true
autorestart=false
startsecs=0
priority=30
stdout_logfile=/var/log/supervisor/migrate.log
stderr_logfile=/var/log/supervisor/migrate_error.log
environment=PYTHONPATH="/app",DATABASE_URL="postgresql://runpod:runpod_secure_pass@localhost:5432/shorts_editor"

# Celery Worker (GPU Tasks) - Priority 40
[program:celery-worker]
command=celery -A core.tasks.celery_app worker --loglevel=info --concurrency=2 --queues=gpu_render,audio_process,video_process --hostname=gpu-worker@%%h --max-tasks-per-child=10
directory=/app
user=root
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/celery_worker.log
stderr_logfile=/var/log/supervisor/celery_worker_error.log
environment=PYTHONPATH="/app",CUDA_VISIBLE_DEVICES="all",CELERY_WORKER_TYPE="gpu"

# Celery Beat (Scheduler) - Priority 41
[program:celery-beat]
command=celery -A core.tasks.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule
directory=/app
user=root
autostart=true
autorestart=true
priority=41
stdout_logfile=/var/log/supervisor/celery_beat.log
stderr_logfile=/var/log/supervisor/celery_beat_error.log
environment=PYTHONPATH="/app"

# Celery Flower (Monitoring) - Priority 42
[program:celery-flower]
command=celery -A core.tasks.celery_app flower --port=5555 --basic_auth=admin:runpod_flower_pass
directory=/app
user=root
autostart=true
autorestart=true
priority=42
stdout_logfile=/var/log/supervisor/celery_flower.log
stderr_logfile=/var/log/supervisor/celery_flower_error.log
environment=PYTHONPATH="/app"

# FastAPI Application Server - Priority 50 (starts last)
[program:fastapi]
command=uvicorn runpod_handler:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info --access-log
directory=/app
user=root
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi_error.log
environment=PYTHONPATH="/app",REDIS_URL="redis://localhost:6379/0",DATABASE_URL="postgresql://runpod:runpod_secure_pass@localhost:5432/shorts_editor"

# Health Monitor - Priority 60 (monitoring service)
[program:health-monitor]
command=python /app/core-svc/monitoring/health_monitor.py
directory=/app
user=root
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/supervisor/health_monitor.log
stderr_logfile=/var/log/supervisor/health_monitor_error.log
environment=PYTHONPATH="/app"

# GPU Monitor - Priority 61 (GPU utilization tracking)
[program:gpu-monitor]
command=python /app/core-svc/monitoring/gpu_monitor.py
directory=/app
user=root
autostart=true
autorestart=true
priority=61
stdout_logfile=/var/log/supervisor/gpu_monitor.log
stderr_logfile=/var/log/supervisor/gpu_monitor_error.log
environment=PYTHONPATH="/app",CUDA_VISIBLE_DEVICES="all"

# Group all services for easy management
[group:database]
programs=postgresql,redis

[group:workers]
programs=celery-worker,celery-beat,celery-flower

[group:api]
programs=fastapi

[group:monitoring]
programs=health-monitor,gpu-monitor